# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jetty]
directory=/root/omh-shims
command=mvn spring-boot:run

[program:mongod]
command=/usr/bin/mongod
